# bin/hrmd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 
HRM server daemon for launchd in MacOS X.
As root:
1) Edit this file to refer to the proper locations of your HRM installation
2) Place the file in /Library/LaunchDaemons/
3) Load it with:
   sudo launchctl load -w /Library/LaunchDaemons/hrmd.plist

Read more at http://developer.apple.com/documentation/Darwin/Reference/ManPages/man5/launchd.plist.5.html

-->
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
        <!-- EDIT THE STRING VALUES OF THESE FOUR VARIABLES: -->
	<dict>
		<key>HRM_HOME</key>
                <!-- HRM installation directory -->
		<string>/Library/WebServer/Documents/hrm</string>
		<key>HRM_DATA</key>
                <!-- Image repository location -->
		<string>/Library/WebServer/Documents/hrm_images</string>
		<key>HRM_SOURCE</key>
                <!-- User raw image subdirectory: -->
		<string>huygens_src</string>
		<key>HRM_DEST</key>
                <!-- User restored image subdirectory: -->
		<string>huygens_dst</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>hrmd</string>
	<key>OnDemand</key>
	<false/>
	<key>ProgramArguments</key>
	<array>
                <!-- The php application is usually here: -->
		<string>/usr/bin/php</string>
		<string>-q</string>
                <!-- EDIT THIS PATH TO HRM: -->
		<string>/Library/WebServer/Documents/hrm/run/runHuygensRemoteManager.php</string>
	</array>

	<key>StandardErrorPath</key>
        <!-- EDIT THIS PATH -->
	<string>/Library/WebServer/Documents/hrm/run/log_error.txt</string>
	<key>StandardOutPath</key>
        <!-- EDIT THIS PATH -->
	<string>/Library/WebServer/Documents/hrm/run/log.txt</string>
	<key>WorkingDirectory</key>
        <!-- EDIT THIS PATH -->
	<string>/Library/WebServer/Documents/hrm/run</string>

	<key>UserName</key>
	<string>root</string>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
